# comexr-fan-control_0.0.1-1_amd64/etc/xdg/autostart/comexr-fan-control.desktop
[Desktop Entry]
Type=Application
Name=comexr-fan-control
Comment=Driver for fans on Clevo machines
Path=/usr/local/bin
Exec=comexr-fan-control
Terminal=false
